# dictation.service
[Unit]
Description=Dictation transcription daemon (faster-whisper)

[Service]
ExecStart=%h/llm-sessions/stt-widget/venv/bin/python %h/llm-sessions/stt-widget/dictation_server.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
